# supervisor/tcpser8.conf
[program:tcpser8]
command=/usr/bin/tcpser -v 7008 -s 115200 -p 8008
priority=1
directory=/dos
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
